# cookbooks/power_led/files/etc/systemd/system/power-led-on.service
[Unit]
Description=Turn on Power LED

[Service]
Type=simple
ExecStart=/usr/local/sbin/power_led_on

[Install]
WantedBy=multi-user.target
